# roles/vendor/redis/templates/redis_cluster.service
[Unit]
Description=Redis key-value database on {{ redis_port }}
After=network.target
[Service]
ExecStart=/usr/local/bin/redis-server redis_{{ redis_port }}.conf --supervised systemd --daemonize no
Type=notify
User=redis
Group=redis
WorkingDirectory=/etc/redis_cluster/redis_{{ redis_port }}
LimitNOFILE=65535
Restart=always

[Install]
WantedBy=multi-user.target
