# roles/vendor/redis/templates/redis.service
[Unit]
Description=Redis key-value database on {{ redis_port }}
After=network.target
[Service]
ExecStart=/usr/bin/redis-server redis_{{ redis_port }}.conf --supervised systemd --daemonize no
ExecStop=/usr/bin/redis-cli -h 127.0.0.1 -p {{ redis_port }} shutdown
Type=notify
User=redis
Group=redis
{% if redis_cluster == false %}
WorkingDirectory=/etc/{{ redis_dir }}
{% else %}
WorkingDirectory=/etc/{{ redis_dir }}/redis_{{ redis_port }}
{% endif %}
LimitNOFILE=65535
Restart=always

[Install]
WantedBy=multi-user.target
